# scripts/aware-bt-connect.service
[Unit]
Description=AWARE Bluetooth speaker auto-connect
After=bluetooth.target bluealsa.service
Before=aware.service
Wants=bluetooth.target bluealsa.service

[Service]
Type=oneshot
User=arduino
Environment=AWARE_BT_SPEAKER_MAC=15:D2:D2:C5:6B:0C
ExecStart=/home/arduino/aware/scripts/connect-bt.sh 10
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
