# scripts/aware.service
[Unit]
Description=AWARE - Autonomous Witness And Response Engine
After=network.target

[Service]
Type=simple
User=aware
WorkingDirectory=/home/aware/aware
ExecStart=/home/aware/aware/.venv/bin/python -m aware.app.main
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
