# scripts/llama-server.service
[Unit]
Description=llama.cpp server for AWARE LLM
After=network.target
Before=aware.service

[Service]
Type=simple
User=arduino
WorkingDirectory=/home/arduino/llama
ExecStart=/home/arduino/llama/llama-b10066/llama-server -m /home/arduino/llama/MiniCPM5-1B-Q4_K_M.gguf --host 127.0.0.1 --port 8080 --reasoning off -c 2048 --mlock --no-mmap
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
